# etc/services.d/rclone/run
#!/usr/bin/with-contenv bash

/usr/bin/rclone ${SUBCMD} ${CONFIG} ${PARAMS}
